# deploy/systemd/discord-bot.service
[Unit]
Description=Discord Bot for Minecraft Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/discord-minecraft-bot
ExecStart=/usr/bin/node /home/pi/discord-minecraft-bot/src/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
